# deps/pi/otto.service
[Unit]
Description=Otto
After=network-online.target
Wants=network-online.target systemd-networkd-wait-online.service
StartLimitIntervalSec=500
StartLimitBurst=5

[Service]
User=otto
Group=otto
Restart=on-failure
RestartSec=5s
WorkingDirectory=/home/otto/ai
ExecStart=pnpm run start:pi

[Install]
WantedBy=default.target